# monitoring/crontab
# Run health monitoring every 15 minutes
*/15 * * * * /usr/local/bin/health-monitor.sh >> /logs/health-monitor.log 2>&1

# Empty line required at end of crontab
